# network-manager.service
[Unit]
Description=Network Manager for Price Checker
After=network.target wpa_supplicant.service
Wants=network.target

[Service]
Type=oneshot
User=root
ExecStart=/opt/price-checker/network-manager.sh auto
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
