# price-checker-kiosk.service
[Unit]
Description=Price Checker Kiosk Mode
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=dietpi
Group=dietpi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/dietpi/.Xauthority
ExecStartPre=/bin/sleep 10
ExecStart=/opt/price-checker/start-kiosk.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
