# price-checker.service
[Unit]
Description=Price Checker Kiosk Server
After=network.target
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=dietpi
Group=dietpi
WorkingDirectory=/opt/price-checker
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=8000

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /tmp /opt/price-checker

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=price-checker

[Install]
WantedBy=multi-user.target

